# deploy/linux/freefsm.service
[Unit]
Description=FreeFSM Field Service Manager
After=network.target postgresql.service

[Service]
Type=simple
User=freefsm
Group=freefsm
EnvironmentFile=-/etc/freefsm.conf
ExecStart=/usr/local/bin/freefsm
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/freefsm /var/log/freefsm
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
